# deploy/revive-thrive.service.template
[Unit]
Description=Revive Thrive Flask App
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/projects/revive-thrive-business-platform
Environment="PATH=/home/ubuntu/projects/revive-thrive-business-platform/venv/bin"
EnvironmentFile=-/home/ubuntu/projects/revive-thrive-business-platform/.env
ExecStart=/home/ubuntu/projects/revive-thrive-business-platform/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
